# cron/roadwarrior-backup.cron

*/5 *	* * *	root	/usr/local/sbin/roadwarrior-backup-trigger.sh

